# build/DEBIAN/postinst
#!/bin/sh
# postinst script for emby-server
#
# see: dh_installdeb(1)

set -e

case "$1" in
  configure)
    getent group emby > /dev/null 2>&1 || groupadd -r emby
    getent passwd emby > /dev/null 2>&1 && usermod -c 'Emby Server' -g emby -d /var/lib/emby -m -s /usr/bin/nologin emby
    getent passwd emby > /dev/null 2>&1 || useradd -r -c 'Emby Server' -g emby -d /var/lib/emby -m -s /usr/bin/nologin emby && passwd -l emby > /dev/null
    getent group video > /dev/null 2>&1 && usermod -aG video emby
    getent group render > /dev/null 2>&1 && usermod -aG render emby
    [ -f /opt/emby-server/lib/ld-2.25.so ] && chmod +x /opt/emby-server/lib/ld-2.25.so
    if [ -x "$(which systemctl)" ]; then
      systemctl daemon-reload
      systemctl enable emby-server || true
      systemctl start emby-server || true
    else
      service emby-server start
    fi
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0

# vim: ts=2 sw=2 et:
